# bannerfilter.cron
#!/bin/sh
/usr/sbin/bannerfilter-update
